# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log

[program:postgresql]
command=/bin/bash -c "exec /usr/lib/postgresql/$(ls /usr/lib/postgresql/)/bin/postgres -D /var/lib/postgresql/data/pgdata -c listen_addresses=127.0.0.1"
user=postgres
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=python /app/app.py
environment=DATABASE_URL="postgresql://guestbook:guestbook123@127.0.0.1:5432/guestbook"
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:apache]
command=apachectl -D FOREGROUND
autorestart=true
priority=30
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
